# nft-tracer.service
[Unit]
Description=NFT Packet Tracer Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nft-tracer-app/backend
ExecStart=/usr/bin/python3 /opt/nft-tracer-app/backend/app.py
Restart=always
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_ENV=production"

# Security
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF CAP_PERFMON
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF CAP_PERFMON
NoNewPrivileges=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nft-tracer

[Install]
WantedBy=multi-user.target
